# config/systemd_backup/jupyter.service
[Unit]
Description=JupyterLab Server
After=network.target

[Service]
Type=simple
ExecStart=/home/nvidia/.venv/bin/jupyter-lab --ip=0.0.0.0 --port=8888 --no-browser
Restart=always
RestartSec=5
WorkingDirectory=/home/nvidia
Environment=PATH=/home/nvidia/.venv/bin:/usr/bin:/bin

[Install]
WantedBy=default.target

